# ld.so.preload
/opt/datadog/./linux-x64/Datadog.Linux.ApiWrapper.x64.so